# openclaw/rootfs/etc/s6-overlay/s6-rc.d/openclaw-init/up
/etc/s6-overlay/scripts/openclaw-init.sh
