# openclaw/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.warning "nginx exited with code ${1}"
