# openclaw/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start nginx (ingress reverse proxy)
# ==============================================================================

# Debug: verify server config exists
if [ ! -f /etc/nginx/servers/ingress.conf ]; then
    bashio::log.error "Missing /etc/nginx/servers/ingress.conf — init may have failed"
    ls -la /etc/nginx/servers/ 2>&1
    exit 1
fi

bashio::log.info "Starting nginx (config: $(wc -l < /etc/nginx/servers/ingress.conf) lines in ingress.conf)..."
exec nginx
