# openclaw/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start nginx (ingress reverse proxy)
# ==============================================================================

bashio::log.info "=== DEBUG: nginx.conf ==="
cat /etc/nginx/nginx.conf
bashio::log.info "=== DEBUG: servers/ ==="
ls -la /etc/nginx/servers/ 2>&1 || echo "servers/ not found"
cat /etc/nginx/servers/*.conf 2>/dev/null || echo "no .conf files"
bashio::log.info "=== DEBUG: listening ports ==="
nginx -T 2>&1 | grep -E "listen|server_name|root|include" || true
bashio::log.info "=== Starting nginx ==="
exec nginx
