# openclaw/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start nginx (ingress reverse proxy)
# ==============================================================================

bashio::log.info "Starting nginx..."
exec nginx -g "daemon off;"
