# openclaw/rootfs/etc/s6-overlay/s6-rc.d/ttyd/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

if ! bashio::config.true 'enable_terminal'; then
    exit 256
fi

bashio::log.warning "ttyd exited with code ${1}"
