# openclaw/rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start ttyd web terminal (conditional)
# ==============================================================================

if ! bashio::config.true 'enable_terminal'; then
    bashio::log.info "Terminal is disabled. Sleeping."
    exec sleep infinity
fi

bashio::log.info "Starting ttyd web terminal..."
exec ttyd \
    --writable \
    --interface 127.0.0.1 \
    --port 7681 \
    --base-path /terminal \
    --cwd /config \
    bash -c 'set -a; source /data/openclaw/env 2>/dev/null; set +a; exec bash'
